# docker-django/container-files/etc/supervisor.d/django.conf
[program:app-uwsgi]
command = uwsgi --ini /opt/app/uwsgi.ini
