# docker-django/container-files/etc/supervisor.d/nginx.conf
[program:nginx]
command = /usr/sbin/nginx
autorestart = false
stderr_logfile = NONE
stdout_logfile = NONE
