# docker-flask/container-files/etc/supervisor.d/flask.conf
[program:app-uwsgi]
command = uwsgi --ini /opt/uwsgi.ini
